# vantage@.service
[Unit]
Description=Vantage Discord Bot (%i)
After=network.target

[Service]
Type=simple
# Replace 'vantage' with the Linux user that runs the bot
User=vantage
# Split-path layout:
#   Code  → /opt/%i/      (git clone)
#   Data  → /var/lib/%i/  (config, cog data, guild files)
#   Logs  → journald (journalctl -u vantage@%i -f)
WorkingDirectory=/opt/%i
Environment=VANTAGE_DATA_DIR=/var/lib/%i
ExecStart=/opt/%i/venv/bin/python launcher.py start
Restart=on-failure
RestartSec=5
# Give the bot 30 seconds to shut down cleanly before systemd kills it
TimeoutStopSec=30
# Log to the journal — view with: journalctl -u vantage@%i -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vantage-%i

[Install]
WantedBy=multi-user.target
