# vprod.service
[Unit]
Description=vprod Discord Bot
After=network.target

[Service]
Type=simple
# System account that runs the bot — created by: sudo python launcher.py system create-user
User=vprodbot
# Layout:
#   Code  -> /opt/vprod/       (git clone, venv)
#   Data  -> /var/lib/vprod/   (config.json, cog_data.json, guilds/, repos/)
#   Token -> /opt/vprod/.env   (DISCORD_TOKEN — chmod 600, never committed)
#   Logs  -> journald           (journalctl -u vprod -f)
WorkingDirectory=/opt/vprod
EnvironmentFile=-/opt/vprod/.env
Environment=VPROD_DATA_DIR=/var/lib/vprod
ExecStart=/opt/vprod/venv/bin/python launcher.py start
Restart=on-failure
RestartSec=5
# Give the bot 30 seconds to shut down cleanly before systemd kills it
TimeoutStopSec=30
# Log to the journal — view with: journalctl -u vprod -f
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vprod

[Install]
WantedBy=multi-user.target
